# infra/files/dagster-code.service
[Unit]
Description=Dagster gRPC Code Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dagster
Group=dagster
WorkingDirectory=/home/dagster/data-platform
EnvironmentFile=/home/dagster/data-platform/.env
ExecStart=/home/dagster/.local/bin/uv run dagster api grpc --host 127.0.0.1 --port 4266 --module-name assets.definitions
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
