# infra/files/dagster-healthcheck.service
[Unit]
Description=Dagster Health Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dagster-healthcheck.sh
